# scripts/com.mem0.gap-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mem0.gap-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/marabook_m1/Documents/APP_HOME/CascadeProjects/windsurf-project/SecondBrain/scripts/monitor_memory_gaps.py</string>
        <string>--alert-only</string>
    </array>

    <key>StartInterval</key>
    <integer>86400</integer> <!-- Run daily (24h) -->

    <key>StandardOutPath</key>
    <string>/tmp/mem0-gap-monitor.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/mem0-gap-monitor-error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
